# dev-util/gerrit-bin/files/gerrit-bin.init2
#!/sbin/openrc-run

depend() {
    need net
    use dns logger mysql postgresql
}

checkconfig() {
    if [ ! -n "$GERRIT_HOME" ] ; then
        eerror "GERRIT_HOME not configured"
        return 1
    fi
    if [ ! -d "$GERRIT_HOME" ] ; then
        checkpath -q -d -o gerrit:gerrit -m 0755 "$GERRIT_HOME" || return 1
    fi
    return 0
}

start() {
    checkconfig || return 1

    JAVA_HOME=$(
           java-config --select-vm=openjdk-bin-11 --jre-home 2> /dev/null \
        || java-config --select-vm=openjdk-bin-17 --jre-home 2> /dev/null )

    COMMAND=$JAVA_HOME/bin/java

    JAVA_PARAMS="$GERRIT_JAVA_OPTIONS -jar $GERRIT_WAR"

    PARAMS="daemon -d $GERRIT_HOME"
    [ -n "$GERRIT_ARGS" ] && PARAMS="$PARAMS $GERRIT_ARGS"

    ebegin "Starting ${SVCNAME}"
    start-stop-daemon --start --quiet --background \
        --make-pidfile --pidfile $GERRIT_PIDFILE \
        --user ${GERRIT_USER} \
        --exec "${COMMAND}" -- $JAVA_PARAMS $PARAMS
    eend $?
}

stop() {
    ebegin "Stopping ${SVCNAME}"
    start-stop-daemon --stop --quiet --pidfile $GERRIT_PIDFILE
    eend $?
}
